# vmagent-process-exporter-supervisor.ini
[program:process-exporter]
command=/data/monitor/exporters/process_exporter/process_exporter --remote.url=http://victoriametrics:8428/api/v1/write --interval=60s --top=10 --label env=production --label region=cn-north
directory=/data/monitor/exporters/process_exporter
user=root
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/process-exporter.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/process-exporter-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PATH="/usr/local/bin:/usr/bin:/bin"
